# debian/pd-mapper.postinst
#!/bin/sh

#DEBHELPER#

rc-update del pd-mapper default || true
rc-update add pd-mapper sysinit || true

exit 0
